# server/crimson-harbor.service
[Unit]
Description=深红港任务公会 API 服务
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/crimson-harbor/server
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target
